# mqtt_tinygs_listen.service
[Unit]
# Legacy single-user example for the original student machine.
# Prefer mqtt_tinygs_listen@.service for new Debian deployments.
Description=TinyGS MQTT listener (latest state JSON)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=student
Group=student
WorkingDirectory=/home/student/synscan_tinygs_tracker
EnvironmentFile=/home/student/synscan_tinygs_tracker/mqtt_tinygs_listen.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/student/synscan_tinygs_tracker/tinygs_mqtt/env/bin/python /home/student/synscan_tinygs_tracker/mqtt_tinygs_listen.py \
  --user ${TINYGS_USER} \
  --station ${TINYGS_STATION} \
  --password ${TINYGS_PASS} \
  --out /home/student/synscan_tinygs_tracker/state.json \
  --frame-topic tinygs/${TINYGS_USER}/${TINYGS_STATION}/cmnd/frame/0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
