# mqtt_tinygs_listen.service
[Unit]
Description=TinyGS MQTT listener (raw + latest JSON)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=student
Group=student
WorkingDirectory=/home/student/01diplomka
EnvironmentFile=/home/student/01diplomka/mqtt_tinygs_listen.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/student/01diplomka/tinygs_mqtt/env/bin/python /home/student/01diplomka/mqtt_tinygs_listen.py \
  --user ${TINYGS_USER} \
  --station ${TINYGS_STATION} \
  --password ${TINYGS_PASS} \
  --out /home/student/01diplomka/state.json \
  --raw-dir /home/student/01diplomka/raw \
  --rx-out /home/student/01diplomka/all_rx.jsonl \
  --frame-topic tinygs/${TINYGS_USER}/${TINYGS_STATION}/cmnd/frame/0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
